# assets/configs/5-selfping.conf
[program:self-ping]
command=/bloop.sh
autorestart=true
